# evermetrics-deploy1/roles/evermetrics/files/upstart/evermetrics-api-auth.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

setuid root
chdir /opt/evermetrics-api

env NODE_ENV="production"
env EVERMETRICS_SECRET="saweeew34-d32d-3ed2-fk56-9cjd8492d843"

exec npm run auth
